# systemd/minio/minio.service
[Unit]
Description=MinIO - S3 Compatible Object Storage
Documentation=https://min.io/docs/minio/
After=network.target

[Service]
Type=simple
User=minio
Group=minio
EnvironmentFile=/etc/default/minio
ExecStart=/usr/local/bin/minio server $MINIO_OPTS $MINIO_VOLUMES
Restart=always
RestartSec=10
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
