# systemd/node-exporter/node-exporter.service
[Unit]
Description=Node Exporter - Host Metrics Exporter
Documentation=https://github.com/prometheus/node_exporter
After=network.target

[Service]
Type=simple
User=node_exporter
Group=node_exporter
ExecStart=/usr/local/bin/node_exporter \
    --path.procfs=/proc \
    --path.sysfs=/sys \
    --path.rootfs=/ \
    --collector.filesystem.mount-points-exclude='^/(sys|proc|dev|host|etc)($$|/)'
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
